# debian/ralt-awesome.postinst
#!/bin/bash

set -e

# pyxtrlock
systemctl enable xtrlock
