# debian/ralt-awesome.postinst
#!/bin/bash

set -e

# pyxtrlock
systemctl enable xtrlock

. /usr/share/debconf/confmodule

db_get ralt_awesome/username

mkdir -p /home/${RET}/.config/
ln -s /usr/share/ralt-awesome/awesome /home/${RET}/.config/awesome
