# debian/ralt-awesome.prerm
#!/bin/bash

set -e

systemctl disable xtrlock
